# supervisord.ini
[supervisord]
user=chrome
nodaemon=true

[program:xvnc]
priority=3000
command=/usr/bin/Xvnc -SecurityTypes None :0
autorestart=true

[program:i3wm]
priority=2000
command=/usr/bin/i3
user=chrome
autorestart=true
environment=HOME="/home/chrome",USER="chrome"

[program:chrome]
priority=1000
command=/usr/bin/chromium-browser --disable-gpu --auto-launch-at-startup --disable-software-rasterizer --disable-dev-shm-usage --window-position=0,0 --no-experiments --noerrdialogs --no-first-run --disable-extensions --disable-file-system
user=chrome
environment=HOME="/home/chrome",USER="chrome"
autorestart=true
startretries=3
stopsignal=HUP
